# splash.service
[Unit]
ConditionKernelCommandLine=splash
ConditionFileNotEmpty=/etc/splash.fb

[Service]
ExecStart=/usr/lib/splash/splash
